# xlpsystem/phabricator /crontab
0 4 * * * /bin/sh -c /root/auto-backup-phabricator.sh